# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:weecloud]
directory=/weecloud-master
command=/usr/bin/node app.js



